# sys-auth/validity-sensor/files/validity-sensor.initd
#!/sbin/openrc-run
# $Id$

depend() {
	use logger
}
name="Validity Fingerprint Service Daemon"
command=/usr/bin/vcsFPService
command_args=""
command_background="false"
