# sys-auth/validity-sensor/files/validity-sensor.initd
#!/sbin/openrc-run
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend() {
	use logger
}
name="Validity Fingerprint Service Daemon"
command=/usr/bin/vcsFPService
command_args=""
command_background="true"
pidfile=${pidfile:-/run/${SVCNAME}.pid}
