# app-laptop/keysid/files/keysid
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

get_event() {
	DEVNODE=`grep -si keyboard /sys/class/input/event*/device/driver/description | grep -o 'event[0-9]\+'`
	if test $? -eq 0; then
		return 0
	else
		return 1
	fi
}

start() {
	ebegin "Starting keysid"

	if get_event; then
		start-stop-daemon --start --quiet --exec /usr/sbin/keysid -- \
		/etc/keysid/keysid.conf
	else
		eerror "Could not gather event interface. Please check your config file"
	fi
	eend $?
}

stop() {
	ebegin "Stopping keysid"
	start-stop-daemon --stop --quiet --exec /usr/sbin/keysid
	eend $?
}
